# infrastructure/systemd/azom-bot.service
[Unit]
Description=Azom dedicated Telegram bot (Jonatan)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=azom
Group=azom
WorkingDirectory=/opt/azom-agent
EnvironmentFile=-/opt/azom-agent/.env
Environment=PYTHONPATH=/opt/azom-agent/skills
Environment=AZOM_CONFIG_DIR=/opt/azom-agent/config
Environment=AZOM_DATA_DIR=/var/lib/azom
Environment=AZOM_USE_MOCK=0
ExecStart=/opt/azom-agent/bin/dedicated-bot.sh
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/azom /var/log/azom /opt/azom-agent/.azom-data

[Install]
WantedBy=multi-user.target
