# infrastructure/systemd/azom-daily-brief.timer
[Unit]
Description=Run Azom daily brief every morning

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
